# grant.conf
[program:grant]
command=/home/da-26/code/grant/run.sh
autostart=false
autorestart=false
startretries=1
user=da-26
stopasgroup=true
killasgroup=true
stopsignal=INT
